# systemd/wifi-hotspot.service
[Unit]
Description=Wi-Fi Hotspot Fallback Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/wifi_hotspot.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
